# systemd/plumb-cursor.service
[Unit]
Description=Plumb — Cursor Agent (A2A Bridge, oneshot)
After=network.target
# Cutover: systemctl stop fang-cursor && systemctl start plumb-cursor

[Service]
Type=simple
WorkingDirectory=/home/plumb
ExecStart=/root/.bun/bin/bun run /home/plumb/src/main.ts wrap "cursor-agent --print" --port 3003 --timeout 300
EnvironmentFile=-/etc/plumb/agent.env
Environment=HOME=/root
Environment=PATH=/root/.local/bin:/root/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
